# systemd/brutedrop.service
[Unit]
Description=BruteDrop
After=sshd.service

[Service]
Type=oneshot
ExecStart=/sbin/brutedrop
StandardOutput=journal
